# eudev/src/rules.d/60-ahci-sss-rescan.rules
#
# Copyright (C) 2026 AuxXxilium <https://github.com/AuxXxilium>
#
# This is free software, licensed under the MIT License.
# See /LICENSE for more information.
#
# Force SCSI host rescan for AHCI controllers that set CAP.SSS (Staggered Spin-Up Support).
#
# When SSS is set the ahci driver defers port scanning, expecting an explicit spin-up trigger.
# Without this rule drives connected to SSS controllers are not visible at boot and require a
# manual driver unbind/rebind to appear.  Affected chipsets include (but are not limited to):
#   ASMedia ASM106x/ASM1064  (1b21:0611/0612/1064)
#   Marvell 88SE9x series    (1b4b:9120/9128/9172/9215/...)
#   JMicron JMB36x           (197b:2360/2362/2363)
#
# Trigger: any PCI device bound to the ahci driver (covers all third-party PCIe SATA cards).
# Action:  write "- - -" to each scsi_host scan file under that PCI device, forcing a full
#          target rescan.  This is equivalent to unbind/rebind but non-disruptive; rescanning
#          a port that already completed its scan is a harmless no-op.
#

ACTION=="add", SUBSYSTEM=="pci", DRIVER=="ahci", \
    RUN+="/bin/sh -c 'for H in /sys/bus/pci/devices/%k/host*/scsi_host/host*/scan; do [ -f \"$H\" ] && echo \"- - -\" > \"$H\"; done'"
